# ubuntu/fvorge/DEBIAN/postinst
chmod 751 /opt/fvorge
cd /opt/fvorge/sbin
ln -s fvorge-lite.pl fvorge
chown -R root.tomcat7 /opt/fvorge
cp /opt/fvorge/etc/init.d/fvorge.sh /etc/init.d/fvorge
update-rc.d fvorge defaults
