# ubuntu/fvorge/DEBIAN/postinst
cd /opt/fvorge/sbin
ln -s fvorge-lite.pl fvorge
ln -s fvorge-manage.pl fvorge-manage
chmod 755 /opt/fvorge
chown -R root.root /opt/fvorge
cp /opt/fvorge/etc/init.d/fvorge.sh /etc/init.d/fvorge
update-rc.d fvorge defaults
