# ubuntu/fvorge/DEBIAN/postrm
update-rc.d -f fvorge remove
rm /etc/init.d/fvorge
