# ubuntu/fvorge/DEBIAN/prerm
update-rc.d -f fvorge remove
rm /opt/fvorge/sbin/fvorge
